# pacman/.config/profile.d/10-pacman.sh
#!/bin/bash
# Path for checkupdate database
export CHECKUPDATES_DB="${XDG_RUNTIME_DIR}/checkupdates-db"

# Search these directories with pacdiff
export DIFFSEARCHPATH='/boot /etc /usr /var'
