# tmux/.config/profile.d/10-tmux.sh
#!/bin/bash
# Set path for tmux' temporay files
export TMUX_TMPDIR="${XDG_RUNTIME_DIR}"
