# xorg/.config/profile.d/10-xauthority.sh
#!/bin/bash
# Set path to Xauthority file
export XAUTHORITY="${XDG_RUNTIME_DIR}/Xauthority"
