# neovim/.config/profile.d/10-nvim.sh
#!/bin/bash
# Set Neovim as the default for multiple applications
export DIFFPROG='nvim -d'
export EDITOR='nvim'
export VISUAL='nvim'
