# libdvdcss/.config/profile.d/10-dvdcss.sh
#!/bin/bash
# Set path to DVDCSS cache
export DVDCSS_CACHE="${XDG_CACHE_HOME}:-${HOME}/.cache}/dvdcss"
